# raxis/systemd/raxis-hygiene.timer
# raxis-hygiene.timer — schedules raxis-hygiene.service every 6h.
#
# OnCalendar fires at 00:00 / 06:00 / 12:00 / 18:00 local time.
# Persistent=true rolls over a missed run on next boot so a host
# that was off at the scheduled hour still gets one sweep
# eventually.
#
# Installation
# ------------
#   cargo xtask hygiene-install-timer --dry-run
#   cargo xtask hygiene-install-timer
#   systemctl --user list-timers raxis-hygiene.timer

[Unit]
Description=Raxis worktree hygiene sweep timer
Documentation=file://__REPO_ROOT__/raxis/guides/operator/18-host-hygiene.md

[Timer]
OnCalendar=*-*-* 00,06,12,18:00:00
Persistent=true
Unit=raxis-hygiene.service

[Install]
WantedBy=timers.target
